# meta-phosphor/common/recipes-connectivity/openipmi/openipmi-2.0.22/ipmi_sim.service
[Unit]
Description=ipmi_sim
After=network.target

[Service]
OOMScoreAdjust=-1000

ExecStart=/usr/bin/ipmi_sim -n
TimeoutSec=120

[Install]
WantedBy=multi-user.target
